# sweettrip.service
[Unit]
Description=SweetTrip Candy Store Application
After=network.target
Wants=network.target

[Service]
Type=forking
User=epfs1
Group=epfs1
WorkingDirectory=/opt/apps/SweetTrip
ExecStart=/opt/apps/SweetTrip/start-sweettrip.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables
Environment=NODE_ENV=development
Environment=PORT=4001

[Install]
WantedBy=multi-user.target


